# system_files/usr/lib/systemd/user/fedora39-toolbox.target
[Unit]
Description=fedora39-toolbox.target for fedora39-toolbox quadlet

[Install]
WantedBy=multi-user.target default.target