# system_files/usr/lib/systemd/user/isengard-cli.target
[Unit]
Description=isengard-cli.target for isengard-cli quadlet

[Install]
WantedBy=default.target

