# system_files/usr/lib/systemd/user/rhel8-ubi-toolbox.target
[Unit]
Description=rhel8-ubi-toolbox.target for rhel8-ubi-toolbox quadlet

[Install]
WantedBy=multi-user.target default.target

