# system_files/usr/lib/systemd/user/rhel9-ubi-toolbox.target
[Unit]
Description=rhel9-ubi-toolbox.target for rhel9-ubi-toolbox quadlet

[Install]
WantedBy=multi-user.target default.target